# chef/site-cookbooks/atlas-seed/files/default/upstart-dnsmasq.conf

description     "DNS/DHCP server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/sbin/dnsmasq -kC /etc/dnsmasq.d/atlas.conf

